# telescope-hole/templates/supervisor.conf
[program:hole]
command={{ hole_dir }}/capture.sh
autostart=true
autorestart=true
startretries=10
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/hole.log

